# openresty.conf
[program:openresty]
command=/opt/openresty/nginx/sbin/nginx -p /opt/openresty/nginx -c /opt/openresty/nginx/conf/nginx.conf  -g 'daemon off;'
autorestart=true